# misc/xod@.service
[Unit]
Description=xod
WantedBy = multi-user.target

[Service]
Type=simple
ExecStart=xod -f /etc/xod/%i -l
KillSignal=SIGINT
Group=audio
